# project/supabase/functions/weekly-ad-refresh/crontab
0 0 * * 0 curl -X POST https://your-project.supabase.co/functions/v1/weekly-ad-refresh